# docker/config/etc/services.d/ossec-logs/run
#!/usr/bin/with-contenv sh

# dumping ossec.log to standard output
exec tail -F /var/ossec/logs/ossec.log
