# docker/config/etc/services.d/filebeat/run
#!/usr/bin/with-contenv sh
echo >&2 "starting Filebeat"

exec /usr/share/filebeat/bin/filebeat -e -c /etc/filebeat/filebeat.yml -path.home /usr/share/filebeat -path.config /etc/filebeat -path.data /var/lib/filebeat -path.logs /var/log/filebeat
