# api/test/integration/env/base/wazuh-manager/sshd.conf
[program:sshd]
directory=/usr/local/
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
redirect_stderr=true

