# ruleset/rules/0860-sysmon_id_13.xml
<!--
  Copyright (C) 2015, Wazuh Inc.
-->

<!--
  Sysmon Event 13 rules ID: 92300 - 92399
-->

<group name="sysmon,sysmon_eid13_detections,windows,">

  <rule id="92300" level="0">
    <if_group>sysmon_event_13</if_group>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)SOFTWARE\\\\(WOW6432NODE\\\\M|M)ICROSOFT\\\\WINDOW(S|S NT)\\\\CURRENTVERSION\\\\(RUN|TERMINAL SERVER\\\\INSTALL\\\\SOFTWARE\\\\MICROSOFT\\\\WINDOWS\\\\CURRENTVERSION\\\\RUN)</field>
    <options>no_full_log</options>
    <description>Added registry content to be executed on next logon</description>
    <mitre>
      <id>T1547.001</id>
    </mitre>
  </rule>

  <rule id="92301" level="12">
    <if_sid>92300</if_sid>
    <field name="win.eventdata.details" type="pcre2">(?i)\.(lnk|vbs|vba)</field>
    <options>no_full_log</options>
    <description>Suspicious file extension detected in registry ASEP to be executed on next logon</description>
    <mitre>
      <id>T1547.001</id>
    </mitre>
  </rule>

  <rule id="92302" level="6">
    <if_sid>92300</if_sid>
    <field name="win.eventdata.image" type="pcre2">(?i)reg\.exe</field>
    <options>no_full_log</options>
    <description>Registry entry to be executed on next logon was modified using command line application reg.exe</description>
    <mitre>
      <id>T1547.001</id>
    </mitre>
  </rule>

  <rule id="92303" level="12">
    <if_sid>92300</if_sid>
    <field name="win.eventdata.details" type="pcre2">(?i)(VNC|tvnserver\.exe)</field>
    <options>no_full_log</options>
    <description>Registry entry to be executed on next logon points to a remote access tool</description>
    <mitre>
      <id>T1547.001</id>
    </mitre>
  </rule>

  <rule id="92304" level="6">
    <if_group>sysmon_event_13</if_group>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)CLASSES\\\\(FOLDER|MS-SETTINGS)\\\\SHELL\\\\OPEN\\\\COMMAND</field>
    <options>no_full_log</options>
    <description>Modified registry key associated to UAC bypass by auto-elevated processes</description>
    <mitre>
      <id>T1548.002</id>
      <id>T1112</id>
    </mitre>
  </rule>

  <rule id="92305" level="12">
    <if_sid>92304</if_sid>
    <field name="win.eventdata.details" type="pcre2">(?i)(cmd|powershell)\.exe</field>
    <options>no_full_log</options>
    <description>Command interpreter added to registry key associated to UAC bypass by auto-elevated processes</description>
    <mitre>
      <id>T1548.002</id>
      <id>T1112</id>
    </mitre>
  </rule>

  <rule id="92306" level="12">
    <if_sid>92304</if_sid>
    <field name="win.eventdata.image" type="pcre2">(?i)(cmd|powershell)\.exe</field>
    <options>no_full_log</options>
    <description>$(win.eventdata.image) added to the registry a subkey associated with UAC bypass by auto-elevated processes</description>
    <mitre>
      <id>T1548.002</id>
      <id>T1112</id>
    </mitre>
  </rule>

  <rule id="92307" level="3">
    <if_group>sysmon_event_13</if_group>
    <field name="win.eventdata.details" type="pcre2">(?i)\.exe</field>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)System\\\\CurrentControlSet\\\\Services</field>
    <options>no_full_log</options>
    <description>Evidence of new service creation found in registry under $(win.eventdata.targetObject) binary is: $(win.eventdata.details)</description>
    <mitre>
      <id>T1543.003</id>
    </mitre>
  </rule>

  <rule id="92308" level="3">
    <if_sid>61615</if_sid>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)CLSID.*LocalServer</field>
    <options>no_full_log</options>
    <description>Possible COM Hijacking evidence found in registry under $(win.eventdata.targetObject) binary is: $(win.eventdata.details)</description>
    <mitre>
      <id>T1546.015</id>
    </mitre>
  </rule>

  <rule id="92309" level="15">
    <if_sid>92308</if_sid>
    <field name="win.eventdata.details" type="pcre2">(?i)AppData</field>
    <options>no_full_log</options>
    <description> COM Hijacking evidence found in registry under $(win.eventdata.targetObject) suspicious binary is: $(win.eventdata.details)</description>
    <mitre>
      <id>T1546.015</id>
    </mitre>
  </rule>

</group>
